# mensacampusnord.service
[Unit]
Description=Beispiel Service

[Service]
ExecStart=/usr/local/bin/MensaCampusNord sub12638_49@46.38.249.177:/openmensa.csicar.de/data.xml